# tailscale/rootfs/etc/s6-overlay/s6-rc.d/setup-extra-ca/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs before tailscaled start
# ==============================================================================

declare ca_certs

# Append ca certificate to /etc/ssl/certs/ca-certificates.crt
if bashio::config.has_value 'ca_cert';
then
    bashio::log.info "Extracting extra Root CA..."
    ca_cert=$(bashio::config "ca_cert")
    echo ${ca_cert} > /usr/local/share/ca-certificates/extra_root_ca.crt
    bashio::log.info "Installing extra Root CA..."
    update-ca-certificates
    bashio::log.info "Installed ca certificate"
else
    bashio::log.notice "Extra Root CA certificate disabled"
fi